# systemd/roon-daemon.service
[Unit]
Description=Roon CLI Daemon
Documentation=https://github.com/EdgarPost/roon-cli
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/roon-daemon
Restart=always
RestartSec=5

# Environment
Environment=NODE_ENV=production

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.config/roon-cli %t
PrivateTmp=true

[Install]
WantedBy=default.target
